# etc/local.d/archroot.stop
#!/usr/bin/env bash
source /etc/archroot.conf

$PRIV umount -Rv $CHROOT/dev
$PRIV umount -Rv $CHROOT/home
$PRIV umount -Rv $CHROOT/lib/modules
$PRIV umount -Rv $CHROOT/proc
$PRIV umount -Rv $CHROOT/run
$PRIV umount -Rv $CHROOT/sys
$PRIV umount -Rv $CHROOT/tmp
$PRIV umount -Rv $CHROOT/var/lib/dbus
 
